# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/etc.conf
[Service]
ExecStart=
ExecStart=systemd-firstboot --prompt-locale --prompt-keymap --prompt-timezone --prompt-root-password --force
ExecStartPre=cp -rp --update=none /usr/share/factory/etc /
